# db_monitor.service
[Unit]
Description=Task Management Database Monitor
After=network.target mysql.service
Requires=mysql.service

[Service]
User=root
WorkingDirectory=/var/task_management/taskmanagement
Environment="PATH=/var/task_management/taskmanagement/venv/bin"
EnvironmentFile=/var/task_management/taskmanagement/.env
ExecStart=/var/task_management/taskmanagement/venv/bin/python3 db_monitor.py
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=yes

[Install]
WantedBy=multi-user.target 