# task_alarm_service.service
[Unit]
Description=Task Alarm Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/task_management_dev
ExecStart=/usr/bin/python3 /var/task_management_dev/task_alarm_service.py
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=task_alarm_service
Environment=PYTHONUNBUFFERED=1
Environment=FIREBASE_CREDENTIALS=%s
Environment=DB_HOST=134.209.149.12
Environment=DB_USER=root
Environment=DB_PASSWORD=123
Environment=DB_NAME=task_db

[Install]
WantedBy=multi-user.target 