# task_management_dev.service
[Unit]
Description=Task Management Development Application
After=network.target mysql.service
Requires=mysql.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/task_management/taskmanagement
Environment="PATH=/var/task_management/taskmanagement/venv/bin"
Environment="FLASK_ENV=development"
Environment="FLASK_DEBUG=1"
EnvironmentFile=/var/task_management/taskmanagement/.env.dev
ExecStart=/var/task_management/taskmanagement/venv/bin/gunicorn --worker-class eventlet -w 1 --bind 0.0.0.0:5001 app:app
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=yes

# Resource limits
LimitNOFILE=65535
TimeoutStartSec=300
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target 